# pr-dashboard.service
[Unit]
Description=PR Dashboard
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/home/pi/Projects/ProgrammingAssistant
ExecStart=/usr/bin/node dist/server/index.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/pi/Projects/ProgrammingAssistant/.env
User=pi

[Install]
WantedBy=multi-user.target
